# target/linux/qualcommax/ipq60xx/base-files/etc/init.d/smp_affinity
#!/bin/sh /etc/rc.common

START=93

enable_affinity_ipq60xx() {
	set_affinity() {
		irq=$(awk "/$1/{ print substr(\$1, 1, length(\$1)-1); exit }" /proc/interrupts)
		[ -n "$irq" ] && echo $2 > /proc/irq/$irq/smp_affinity
	}
	
	# cpus available, 1,2,4,8, shit is a quad core.

	# affine ce interrupts to cpu1
	set_affinity 'ce0' 2
	set_affinity 'ce1' 2
	set_affinity 'ce2' 2
	set_affinity 'ce3' 2
	set_affinity 'ce5' 2
	set_affinity 'ce7' 2
	set_affinity 'ce9' 2
	set_affinity 'ce10' 2
	set_affinity 'ce11' 2
       
       	# assign lan/wan to cpu2
	set_affinity 'edma_txcmpl' 4
	set_affinity 'edma_rxfill' 4
	set_affinity 'edma_rxdesc' 4
	set_affinity 'edma_misc' 4
	
        # then we copy the one from ipq807x
	# assign 4 rx interrupts to each core
	set_affinity 'reo2host-destination-ring1' 1
	set_affinity 'reo2host-destination-ring2' 2
	set_affinity 'reo2host-destination-ring3' 4
	set_affinity 'reo2host-destination-ring4' 8

	# assign 3 tcl completions to last 3 CPUs
	set_affinity 'wbm2host-tx-completions-ring1' 2
	set_affinity 'wbm2host-tx-completions-ring2' 4
	set_affinity 'wbm2host-tx-completions-ring3' 8

	# assign 3 ppdu mac interrupts to last 3 cores
	set_affinity 'ppdu-end-interrupts-mac1' 2
	set_affinity 'ppdu-end-interrupts-mac2' 4
	set_affinity 'ppdu-end-interrupts-mac3' 8
	
	# affine usb to last core
	set_affinity 'xhci-hcd:usb1' 8

}

boot() {
	enable_affinity_ipq60xx
}
